# ghost/rootfs/etc/services.d/ghost/run
#!/usr/bin/with-contenv sh
cd /ghost
echo "[*] Starting Ghost"
exec s6-setuidgid ghost /usr/local/bin/node index
